# dist/systemd/podspawn-cleanup.service
[Unit]
Description=Podspawn cleanup pass (single run)
Documentation=https://podspawn.dev/docs/cli/admin-commands
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/podspawn cleanup
